# systemd/lmstudio-secondary.service
[Unit]
Description=LM Studio Headless Service (Secondary, Isolated)
# Start only after network is ready
After=network-online.target
Wants=network-online.target
# This service conflicts with the original and the primary service
Conflicts=lmstudio.service lmstudio-primary.service

[Service]
# Run as your specific user, not root
# IMPORTANT: Replace <YOUR_USERNAME> with your actual Linux username
User=<YOUR_USERNAME>
Group=<YOUR_USERNAME>

# The script that runs the application with the specific argument
ExecStart=/home/sterl/Downloads/Headless-LMStudio-Server-main/scripts/lms-launcher.sh --secondary-headless

# Restart automatically if it crashes
Restart=on-failure
RestartSec=10

# Ensure service stops cleanly without requiring manual force kill
TimeoutStopSec=10

# Set timeout for startup
TimeoutStartSec=300

[Install]
# Starts at boot for the specified user
WantedBy=multi-user.target
